# pgbouncer/files/supervisor.conf
{%- set service = salt['pillar.get']('supervisor:server:service:'+service_name) %}

[program:pgbouncer]
command=pgbouncer /etc/pgbouncer/pgbouncer.ini --user pgbouncer
stdout_logfile=/var/log/pgbouncer/pgbouncer.log
stderr_logfile=/var/log/pgbouncer/pgbouncer_error.log
user=pgbouncer
autostart=true
autorestart=true
